# dseuhid.service
[Unit]
Description=DualSense UHID Proxy
After=systemd-udevd.service

[Service]
Type=simple
RuntimeDirectory=dseuhid
RuntimeDirectoryMode=0755
ExecStart=/usr/bin/dseuhid
Environment=RUST_LOG=info
ExecReload=/usr/bin/edgemap reload
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
